# docker-nginx.service
[Unit]
Description=Nginx container environment
Requires=docker.service
After=docker.service

[Service]
Restart=always
Type=forking
ExecStart=/bin/bash /etc/docker-nginx/docker-nginx.sh start
ExecStop=/bin/bash /etc/docker-nginx/docker-nginx.sh stop

[Install]
WantedBy=default.target
